# scripts/supervisor.conf
[program:gearman-worker]
command=php /var/www/SERVERNAME/vendor/nrwtaylor/stack-agent-thing/src/worker.php
autostart=true
autorestart=false
numprocs=3
process_name=gearman-worker-%(process_num)s
